# configs/qbittorrent.openrc
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md

command="/usr/local/bin/qbittorrent-nox"
command_args="--profile=/usr/local/etc/"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after sshd firewall
}
